# bsp/arch-iMX6/apps/base/etc/inittab
::sysinit:/etc/rcS
::restart:/sbin/init
ttymxc0::askfirst:-/bin/sh
